# systemd/comunito-portal.service
[Unit]
Description=Comunito Portal
After=network-online.target tailscaled.service
Wants=network-online.target tailscaled.service
StartLimitIntervalSec=0
TimeoutStartSec=0

[Service]
User=pi
WorkingDirectory=/home/pi/comunito_portal
ExecStart=/home/pi/comunito_portal/comunito-venv/bin/python -u portal_alpr_full.py
Restart=always
RestartSec=5

# Watchdog: systemd mata y reinicia si el proceso no pinga en 120s
WatchdogSec=120
NotifyAccess=all

# Limitar RAM — OOM killer elimina ESTE proceso antes de congelar el sistema
MemoryMax=700M
MemorySwapMax=0
OOMPolicy=kill
OOMScoreAdjust=500

# TimeoutStartSec=120 permite que ALPR model loading tome hasta 2min
# Si tarda más, systemd mata y reinicia (ALPR retry loop lo recupera)
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
